# root/var/spool/cron/crontabs/mirror
0 */4 * * * nvd-update >> /var/log/cron.log 2>&1
